# plex/archive/rootfs/etc/services.d/plex/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Plex Media Server
# Runs the Plex Media Server
# ==============================================================================
bashio::log.info 'Starting the Plex Media Server...'
library_path=$(bashio::config 'library_path')

export LD_LIBRARY_PATH=$(bashio::config 'library_path')
export PLEX_MEDIA_SERVER_HOME=$LD_LIBRARY_PATH
export PLEX_MEDIA_SERVER_MAX_PLUGIN_PROCS=6
export PLEX_MEDIA_SERVER_INFO_DEVICE="Hass.io"
export PLEX_MEDIA_SERVER_APPLICATION_SUPPORT_DIR=/data

# Run the Plex Media Server
exec /usr/lib/plexmediaserver/Plex\ Media\ Server
